# com.keboola.triage-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keboola.triage-automation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eliskasvobodova/github/Linear Product Feedback/run_triage.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday at 7:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Tuesday at 7:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Wednesday at 7:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Thursday at 7:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Friday at 7:00 AM -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/eliskasvobodova/github/Linear Product Feedback/triage-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eliskasvobodova/github/Linear Product Feedback/triage-launchd-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
